# supervisor.ini
[program:traveller-lite]
directory = /traveller-lite
command = venv/bin/gunicorn
autostart = true
autorestart = true
startretries = 3