# supervisor.ini
[program:traveller-lite]
directory = /traveller-lite
command = gunicorn
autostart = true
autorestart = true
startretries = 3


[program:emailer]
directory = /traveller-lite
command = gunicorn
autostart = true
autorestart = true
startretries = 3
