# supervisor.ini
[program:traveller-lite]
directory = /traveller-lite
command = gunicorn
autostart = true
autorestart = true
startretries = 3


[program:huey]
directory = /traveller-lite
command = huey_consumer app.huey.run
autostart = true
autorestart = true
startretries = 3